# named-root.service
[Unit]
Description=Reset named as ID 0

[Service]
Type=oneshot

ExecStart=/usr/libexec/named-root.sh
